# src/cronjob.cron
0 5 * * * /usr/local/bin/python /source/class_reminders.py >> /var/log/cron.log 2>&1
